# web/codes/supervisor.conf
[supervisord]
nodaemon = true

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --emperor /etc/uwsgi/vassals

[program:nginx-app]
command = service nginx start
